# sa-policy.cron
MAILTO=""
* * * * * root /opt/sa-policy/sa-policy.pl & >/dev/null 2>&1
